# root/etc/services.d/uptime-kuma/run
#!/bin/sh
export HOME=/config
export DATA_DIR=/config
export NODE_ENV=production

cd /app/uptime-kuma
exec s6-setuidgid bsd node server/server.js
